# scripts/systemd/oracle-iteration.timer
[Unit]
Description=Run Oracle Lab agent iteration daily

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
